# etc/supervisor_available/bedrock.conf
[program:bedrock]
command = /app/docker/run-prod.sh
numprocs = 1
autostart = true
user = webdev
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
