# etc/supervisor_available/cron_db.conf
[program:cron_db]
command = /app/bin/run-clock.sh
numprocs = 1
autostart = true
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
