# etc/supervisor_available/cron.conf
[program:cron]
command = /app/bin/cron.py
numprocs = 1
autostart = true
user = webdev
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
