# etc/supervisor_available/cron_l10n.conf
[program:cron_l10n]
command = /app/bin/cron_l10n.py
numprocs = 1
autostart = true
user = webdev
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
