# etc/cron.d/bedrock-stage.tmpl
#
# {{ header }}
#

MAILTO="webops-cron@mozilla.com,cron-bedrock@mozilla.com"

# add date/time info to cron log
0    * * * * {{ user }} date {{ log }}

*/10 * * * * {{ user }} {{ source }}/bin/update-scripts/stage/update-stage-php.sh {{ log }}
*/15 * * * * {{ user }} {{ source }}/bin/update-scripts/stage/update-stage-locale.sh {{ log }}

# bug 1142850
0    * * * * {{ user }} {{ source }}/bin/update-scripts/stage/update-externalfiles.sh {{ log }}

*/5  * * * * {{ django_manage }} rnasync {{ log }}

# bug 996144
1 */2 * * * {{ django_cron }} update_tweets {{ log }}

# bug 1072999
# Have to do this in the SRC dir because rsync destroys git repos in WWW dir.
41 * * * * {{ django_src_manage }} update_security_advisories --quiet {{ log }}

# bug 1087533
42 * * * * {{ django_cron }} update_reps_ical {{ log }}

# bug 1128587
38 * * * * {{ django_manage }} runscript update_firefox_os_feeds {{ log }}

# bug 1116511
# every tuesday midnight
0 0 * * 2 {{ django_manage }} runscript update_tableau_data {{ log }}

