# etc/cron.d/bedrock-stage
MAILTO="webops-cron@mozilla.com,cron-bedrock@mozilla.com"
BEDROCK=/data/bedrock-stage/src/www.allizom.org-django/bedrock
BEDROCK_DEPLOY=/data/bedrock-stage/deploy
# TODO: update rest of cron jobs to use the above env variables

# webdev requests that prod auto-updates the PHP side, for similar to old
0-59/10 * * * * root /data/bedrock-stage/src/update-www.allizom.org.sh

*/15 * * * * root /data/bedrock-stage/src/update-www.allizom.org-svn-locale.sh
# outputs when the run completed to http://www.allizom.org/media/locale_finished.txt

# In testing for bug 753566 resolution...
# not completely sure how I want to solve this yet
#
# run the crons from the admin node, but on the web nodes... does extra work, but works
0 * * * * *  root echo "cd /data/www/www.allizom.org-django/bedrock; python manage.py cron update_feeds | /usr/bin/issue-multi-command bedrock-stage
*/5    * * * * root cd /data/bedrock-stage/src/www.allizom.org-django/bedrock; python manage.py rnasync > /dev/null 2>&1

# bug 996144
1 */2 * * * root cd /data/bedrock-stage/src/www.allizom.org-django/bedrock; python manage.py cron update_tweets > /dev/null 2>&1

#
# alternative- run the cron twice right on the admin node, passing --settings=settings/memcache-XXX.py... will this work?
# requires maintaining 2 extra settings files, at least wrt memcache settings
#0 * * * * * root cd /data/www.allizom.org-django/bedrock; python manage.py --settings=settings/memcache-phx1.py cron update_feeds 2>&1 1> /dev/null | grep -v 'old-style Playdoh layout'"
#0 * * * * * root cd /data/www.allizom.org-django/bedrock; python manage.py --settings=settings/memcache-scl3.py cron update_feeds 2>&1 1> /dev/null | grep -v 'old-style Playdoh layout'"

# bug 1072999
41 * * * * root cd /data/bedrock-stage/src/www.allizom.org-django/bedrock; python manage.py update_security_advisories --quiet > /dev/null 2>&1

# bug 1087533
# TODO: redirect output to null after initial deployment and verification
42 * * * * root cd $BEDROCK; python manage.py cron update_reps_ical
43 * * * * root cd $BEDROCK; python manage.py cron update_openstandard_feeds && $BEDROCK_DEPLOY && python manage.py cron update_openstandard_on_homepage
